# deploy/sd/zeitspiegel-firstboot.service
# Installed by firstrun.sh on the first boot; runs the staged provisioning
# from the SD card's boot partition on every boot until it succeeds (needs
# internet once, e.g. ethernet to your router), then disables itself.
[Unit]
Description=Zeitspiegel first-boot provisioning (from /boot/firmware/zeitspiegel)
Wants=network-online.target
After=network-online.target
ConditionPathExists=!/etc/zeitspiegel/.provisioned

[Service]
Type=oneshot
TimeoutStartSec=45min
ExecStart=/bin/bash /boot/firmware/zeitspiegel/firstboot.sh

[Install]
WantedBy=multi-user.target
