# deploy/zeitspiegel.service
# deploy/zeitspiegel.service — appliance unit (NFR-5, NFR-9).
# Installed by setup.sh to /etc/systemd/system/zeitspiegel.service.

[Unit]
Description=Zeitspiegel time-delayed video mirror
Documentation=https://github.com/danielmschmidt/zeitspiegel/tree/main/docs
# Start the mirror as early as possible — it needs no network (FR-12, ≤25 s
# to fullscreen). The HTTP server binds :80 immediately; clients reach it
# once they join the appliance's Wi-Fi, no ordering against the network
# needed. Take ownership of the HDMI console so no getty shares it.
After=local-fs.target
Conflicts=getty@tty1.service

[Service]
ExecStart=/usr/local/bin/zeitspiegel --config /etc/zeitspiegel/config.toml
Restart=always
RestartSec=2
# tmpfs working dir at /run/zeitspiegel; exported clips land here and
# vanish on power loss — no persistent writes on the read-only root (NFR-9).
RuntimeDirectory=zeitspiegel
# Needs direct DRM/KMS (SDL KMSDRM, no X/Wayland) plus /dev/video0. Either
# root or a user in the video+render groups works; we keep User=root for v1
# simplicity on a sealed single-purpose appliance (read-only root, LAN-only).
User=root
# stdout/stderr go to the journal by default; logs stay volatile (NFR-9).

[Install]
WantedBy=multi-user.target
