# deploy/zeitspiegel.service
# deploy/zeitspiegel.service — appliance unit (NFR-5, NFR-9).
# Installed by setup.sh to /etc/systemd/system/zeitspiegel.service.

[Unit]
Description=Zeitspiegel time-delayed video mirror
Documentation=https://github.com/danielmschmidt/zeitspiegel/tree/main/docs
# Ordered after the network so the web UI binds cleanly when Wi-Fi is up,
# but deliberately NOT Requires= — the mirror must start with Wi-Fi down
# (docs/DEPLOYMENT.md: the web UI appears when the network does).
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/zeitspiegel --config /etc/zeitspiegel/config.toml
Restart=always
RestartSec=2
# tmpfs working dir at /run/zeitspiegel; exported clips land here and
# vanish on power loss — no persistent writes on the read-only root (NFR-9).
RuntimeDirectory=zeitspiegel
# Needs direct DRM/KMS (SDL KMSDRM, no X/Wayland) plus /dev/video0. Either
# root or a user in the video+render groups works; we keep User=root for v1
# simplicity on a sealed single-purpose appliance (read-only root, LAN-only).
User=root
# stdout/stderr go to the journal by default; logs stay volatile (NFR-9).

[Install]
WantedBy=multi-user.target
